# systemd-mithlond.service
[Unit]
Description=Mithlond

[Service]
Type=oneshot
ExecStart=/lib/systemd/mithlond

[Install]
WantedBy=multi-user.target
